# config/linnea.service
# Linnea systemd unit. One-time install:
#   make && sudo make install
#   sudo cp config/linnea.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now linnea
# Routine deploy after a rebuild — zero downtime, no dropped connections:
#   make && sudo make install && sudo systemctl reload linnea
# reload sends SIGUSR2: the master re-execs the new binary in place (same
# PID), the new generation adopts the still-open listeners and starts
# serving, then the old workers drain. If the new binary rejects the
# config the reload is refused and the old generation keeps serving. A
# reload cannot change the listener set (host/port); for that use restart.
#
# kTLS needs the kernel tls module; make it persistent across boots with:
#   echo tls | sudo tee /etc/modules-load.d/ktls.conf
#
# systemd supervises the master; the master supervises the workers.
# Stopping the unit SIGTERMs the whole control group: workers drain —
# stop accepting, finish open connections, then exit — and systemd
# waits for the cgroup (up to TimeoutStopSec, then SIGKILL). The
# workers' PR_SET_PDEATHSIG covers a master that dies on its own.

[Unit]
Description=Linnea HTTP/TLS server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=linnea
Group=linnea
# the config's root and log paths are relative to this directory
WorkingDirectory=/home/linnea/linnea
ExecStart=/usr/local/bin/linnea /home/linnea/linnea/config/linnea-tls.json
# reload = zero-downtime binary upgrade (SIGUSR2, re-exec in place). The
# PID is preserved across the exec, so MainPID tracking survives.
ExecReload=/bin/kill -USR2 $MAINPID
Restart=on-failure
RestartSec=1

# Bind 80/443 without root and without file capabilities on the binary —
# a rebuilt bin/linnea keeps working with no setcap re-run.
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
NoNewPrivileges=true

# Hardening: read-only /usr /boot /etc, private /tmp, no kernel tunables.
# /home stays visible: the binary, config, docroot, log and certs live there.
ProtectSystem=full
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
